# deploy/sumeru-v2.service
[Unit]
Description=Sumeru v2 Host
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/repos/sumeru
Environment=PATH=%h/.local/share/npm/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=SUMERU_PORT=7900
EnvironmentFile=-%h/.config/sumeru/env
ExecStart=/usr/bin/node packages/host/dist/main.js %h/repos/sumeru
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
